# jenkins-approval-bot.service
[Unit]
Description=Jenkins Approval Bot - Telegram机器人审批服务
Documentation=https://github.com/your-repo/jenkins-approval-bot
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
# 服务类型
Type=simple
Restart=on-failure
RestartSec=10s

# 用户和工作目录
User=jenkins-bot
Group=jenkins-bot
WorkingDirectory=/opt/jenkins-approval-bot

# 环境变量
Environment=PATH=/opt/jenkins-approval-bot/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/jenkins-approval-bot
Environment=PYTHONUNBUFFERED=1

# 启动命令
ExecStart=/opt/jenkins-approval-bot/venv/bin/python app.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/jenkins-approval-bot/logs /opt/jenkins-approval-bot/config
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 系统调用限制
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jenkins-approval-bot

[Install]
WantedBy=multi-user.target
